# templates/default/init_script.erb
#!/bin/sh
#
# <%= @application_name %>_<%= @rails_env %> - this script starts and stops rails application
#
# chkconfig:   - 85 15
#
# config:      <%= File.join(@root_folder, '/shared/chef_unicorn.rb') %>
# pidfile:     <%= File.join(@root_folder, '/shared/pids/unicorn.pid') %>

# Source function library.
. /etc/rc.d/init.d/functions

# Source networking configuration.
. /etc/sysconfig/network

# Check that networking is up.
[ "$NETWORKING" = "no" ] && exit 0

set -u
set -e

# Change these to match your app:
APP_NAME="<%= @application_name %>"
APP_ROOT="<%= File.join(@root_folder, 'current') %>"
PID="<%= File.join(@root_folder, 'shared/pids/unicorn.pid') %>"
ENV="<%= @rails_env %>"

UNICORN_OPTS="-D -E $ENV -c <%= File.join(@root_folder, 'shared/config/chef_unicorn.rb') %>"

SET_PATH="cd $APP_ROOT"
CMD="$SET_PATH; ../../shared/scripts/rvm_wrapper.sh bundle exec unicorn $UNICORN_OPTS"

old_pid="$PID.oldbin"

cd $APP_ROOT || exit 1

sig () {
    test -s "$PID" && kill -$1 `cat $PID`
}

oldsig () {
    test -s $old_pid && kill -$1 `cat $old_pid`
}

start_unicorn () {
    if [ `whoami` != <%= @user %> ]; then
        su - <%= @user %> -c "$CMD"
    else
        bash -c "$CMD"
    fi
}

case ${1-help} in
    start)
        sig 0 && echo >&2 "Already running" && exit 0
	start_unicorn
        ;;
    stop)
        sig QUIT && exit 0
        echo >&2 "Not running"
        ;;
    kill)
        sig KILL && exit 0
        echo >&2 "Not running"
        ;;
    force-stop)
        sig TERM && exit 0
        echo >&2 "Not running"
        ;;
    restart)
        if [ -s "$PID" ]; then
	    if [ -n "$(ps -p `cat $PID`|grep -v PID)" ]; then
		kill -TERM `cat $PID`
	    else
		echo "Process does not exist. We just delete pid file and start unicorn"
		rm -f $PID
	    fi
	fi
	    start_unicorn
        ;;
    reload)
	if [ -s "$PID" ]; then
            kill -USR2 `cat $PID`
	else
	    start_unicorn
	fi
        ;;
    upgrade)
        sig USR2 && exit 0
        echo >&2 "Couldn't upgrade, starting '$CMD' instead"
	start_unicorn
        ;;
    rotate)
        sig USR1 && echo rotated logs OK && exit 0
        echo >&2 "Couldn't rotate logs" && exit 1
        ;;
    *)
        echo >&2 "Usage: $0 <start|stop|restart|upgrade|rotate|force-stop>"
        exit 1
        ;;
esac
